# src/systemd/pygitdatback-noui.timer
[Unit]
Description=Run PyGitDatBack No-UI Service Weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target